# helpers/system_statistics_server/sss_location.timer
[Unit]
Description=Runs sss_location.py every day.

[Timer]
OnBootSec=0min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target