# helpers/system_statistics_server/sss_power.timer
[Unit]
Description=Runs sss_power.py every minute.

[Timer]
OnUnitActiveSec=60s
OnBootSec=0min

[Install]
WantedBy=timers.target