# helpers/configurations_update_server/cus.timer
[Unit]
Description=Starts configurations update server on boot.

[Timer]
OnBootSec=0min

[Install]
WantedBy=timers.target